# upstart/msiof-worker.conf
description "MSIOF Worker"
author "Ashley Hindle"

start on startup
stop on shutdown
respawn

exec sudo -u www-data /git/msiof/worker/worker.php
