# upstart/msiof-worker.conf
description "MSIOF Worker"
author "Ashley Hindle"

start on startup
stop on shutdown
respawn

exec sudo -u msiof-worker /etc/msiof/worker
